# templates/default/service_etcd-containerized.service.erb
[Unit]
Description=The Etcd Server container
After=docker.service
Requires=docker.service
PartOf=docker.service

[Service]
EnvironmentFile=<%= node['is_apaas_openshift_cookbook']['etcd_conf_dir'] %>/etcd.conf
ExecStartPre=-/usr/bin/docker rm -f <%= node['is_apaas_openshift_cookbook']['etcd_service_name'] %>
ExecStart=/usr/bin/docker run --name <%= node['is_apaas_openshift_cookbook']['etcd_service_name'] %> --rm -v /var/lib/etcd:/var/lib/etcd:z -v <%= node['is_apaas_openshift_cookbook']['etcd_conf_dir'] %>:<%= node['is_apaas_openshift_cookbook']['etcd_conf_dir'] %>:ro --env-file=<%= node['is_apaas_openshift_cookbook']['etcd_conf_dir'] %>/etcd.conf --net=host --entrypoint=/usr/bin/etcd <%= node['is_apaas_openshift_cookbook']['openshift_docker_etcd_image'] %>
ExecStop=/usr/bin/docker stop <%= node['is_apaas_openshift_cookbook']['etcd_service_name'] %>
SyslogIdentifier=<%= node['is_apaas_openshift_cookbook']['etcd_service_name'] %>
Restart=always
RestartSec=5s

[Install]
WantedBy=docker.service

