# templates/default/service_master-api-containerized.service.erb
[Unit]
Description=Atomic OpenShift Master API
Documentation=https://github.com/openshift/origin
After=etcd_container.service
Wants=etcd_container.service
Before=atomic-openshift-node.service
After=docker.service
PartOf=docker.service
Requires=docker.service

[Service]
EnvironmentFile=/etc/sysconfig/atomic-openshift-master-api
Environment=GOTRACEBACK=crash
ExecStartPre=-/usr/bin/docker rm -f atomic-openshift-master-api
ExecStart=/usr/bin/docker run --rm --privileged --net=host --name atomic-openshift-master-api --env-file=/etc/sysconfig/atomic-openshift-master-api -v <%= node['is_apaas_openshift_cookbook']['openshift_data_dir'] %>:<%= node['is_apaas_openshift_cookbook']['openshift_data_dir'] %> -v /var/log:/var/log -v /var/run/docker.sock:/var/run/docker.sock -v <%= node['is_apaas_openshift_cookbook']['openshift_common_base_dir'] %>:<%= node['is_apaas_openshift_cookbook']['openshift_common_base_dir'] %> -v /etc/pki:/etc/pki:ro <%= node['is_apaas_openshift_cookbook']['openshift_docker_master_image'] %>:${IMAGE_VERSION} start master api --config=${CONFIG_FILE} $OPTIONS
ExecStartPost=/usr/bin/sleep 10
ExecStop=/usr/bin/docker stop atomic-openshift-master-api
LimitNOFILE=131072
LimitCORE=infinity
WorkingDirectory=<%= node['is_apaas_openshift_cookbook']['openshift_data_dir'] %>
SyslogIdentifier=atomic-openshift-master-api
Restart=always
RestartSec=5s

[Install]
WantedBy=docker.service
WantedBy=atomic-openshift-node.service
